# build/DEBIAN/postinst
#!/bin/sh
set -e

SERVICE="power-manager"
echo "enable service..."
systemctl enable $SERVICE
echo "start service..."
systemctl start $SERVICE

