# PowerManager/DEBIAN/postinst
#!/bin/sh
set -e
deb-systemd-helper enable /usr/lib/systemd/system/org.sensey.PowerManager.service
deb-systemd-helper start org.sensey.PowerManager