# PowerManager/DEBIAN/postinst
#!/bin/sh
set -e

SERVICE="org.sensey.PowerManager"
echo "enable service..."
systemctl enable $SERVICE
echo "start service..."
systemctl start $SERVICE

