# PowerManager/DEBIAN/postinst
#!/bin/sh
set -e
systemctl enable /usr/lib/systemd/system/org.sensey.PowerManager.service
systemctl start org.sensey.PowerManager