# PowerManager/DEBIAN/prerm
#!/bin/sh
set -e
deb-systemd-helper stop org.sensey.PowerManager
deb-systemd-helper disable /usr/lib/systemd/system/org.sensey.PowerManager.service
