# PowerManager/etc/init/power-manager.conf
# powermanager - powermanager daemon

description     "powermanager daemon"

start on started lightdm
stop on stopping lightdm


expect fork
respawn

exec /usr/bin/power-manager